# supervisord.conf
[supervisord]
logfile=/tmp/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10) 
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false            ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200

[program:django]
command=/usr/local/python278/bin/python manage.py runserver 0.0.0.0:8000
directory=/root/apps/prj01
stopasgroup=true

[program:celeryd]
command=/usr/local/python278/bin/python manage.py celeryd -l info
directory=/root/apps/prj01
stopasgroup=true
